# install/supervisor/HttpGateway.conf
[program:HpcGateway]
command=dotnet MicroZero.Hpc.Gateway.dll
directory=/usr/zero/gateway
user=root
stopsignal=INT
autostart=true
autorestart=true
startsecs=1
stderr_logfile=/usr/zero/logs/gateway.err.log
stdout_logfile=/usr/zero/logs/gateway.out.log
