# install/supervisor/WebMonitor.conf
[program:ZeroMonitor]
command=dotnet ZeroNet.Devops.Monitor.dll
directory=/usr/zero/monitor
environment=ASPNETCORE__ENVIRONMENT=Production
user=root
stopsignal=INT
autostart=true
autorestart=true
startsecs=1
stderr_logfile=/usr/zero/logs/monitor.err.log
stdout_logfile=/usr/zero/logs/monitor.out.log