# install/supervisor/conf.d/AuthorityCenter.conf
[program:AuthorityCenter]
command=dotnet AuthorityCenter.Api.dll
directory=/root/zero/AuthorityCenter
environment=ASPNETCORE__ENVIRONMENT=Production
user=root
stopsignal=INT
autostart=true
autorestart=true
startsecs=1
stderr_logfile=/root/zero/logs/AuthorityCenter/err.log
stdout_logfile=/root/zero/logs/AuthorityCenter/out.log
