# install/supervisor/conf.d/Background.conf
[program:Background]
command=dotnet ManageApplication.dll
directory=/root/zero/manage
user=root
stopsignal=INT
autostart=true
autorestart=true
startsecs=1
stderr_logfile=/root/zero/logs/Background/err.log
stdout_logfile=/root/zero/logs/Background/out.log