# install/supervisor/conf.d/SimpleHis.conf
[program:SimpleHIS]
command=dotnet SimpleHis.All.dll 
directory=/root/zero/his
user=root
stopsignal=INT
autostart=true
autorestart=true
startsecs=1
stderr_logfile=/root/zero/logs/SHIS/err.log
stdout_logfile=/root/zero/logs/SHIS/out.log